# docker/rootfs/etc/s6-overlay/s6-rc.d/xvfb/run
#!/bin/bash
exec Xvfb :99 -screen 0 ${XVFB_RESOLUTION} -ac +extension GLX +render -noreset
