# docker/rootfs/etc/s6-overlay/s6-rc.d/fluxbox/run
#!/command/with-contenv bash
export DISPLAY=:99
exec fluxbox
