# docker/rootfs/etc/s6-overlay/s6-rc.d/fluxbox/run
#!/bin/bash
export DISPLAY=:99
exec fluxbox
